# 50-usb-disk.rules
ACTION=="add", KERNEL=="sd*", SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", RUN+="/bin/systemctl --no-block start usb-disk-copy-@%k.service"